# crontab
#START

58 4 * * * /app/clean_mem.sh

0 5 * * * cd /app && /usr/bin/python3 ms_rewards_farmer.py --fast --headless --telegram "${TOKEN}" "${CHAT_ID}" > /proc/1/fd/1 2>/proc/1/fd/2

#END
